# docker/supervisor/queue.ini
[program:my-app]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www/html/my-app

command=php artisan queue:work --queue=high,default --sleep=1 --timeout=3600
autostart=true
autorestart=true  
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/html/my-app/storage/logs/worker.log
stopwaitsecs=3600
stdout_logfile_maxbytes=5MB 
